# iso-generator/root/files/ansible/iso/files/repository/temp.service
[Unit]
Description=Installation Service

[Service]
ExecStart=/root/deploy.sh

[Install]
WantedBy=multi-user.target
